# detections/sigma/custom/EDRGap_WMI_Permanent_Subscription.yml
title: WMI Permanent Event Subscription Created with Suspicious Consumer (UH_WMI Persistence)
id: e4f5a6b7-c8d9-0123-efab-444555666777
status: experimental
description: |
    Detects creation of WMI permanent event subscriptions where the consumer
    command line references a known attacker tool or interpreter.

    WMI permanent subscriptions (__EventFilter + __EventConsumer +
    __FilterToConsumerBinding) are used for fileless persistence that survives
    reboots and executes arbitrary commands entirely inside the WMI subsystem,
    without writing a scheduled task or registry run key.

    IWbemServices::PutInstance is the COM interface used to create the
    subscription. It has no named DLL import and cannot be detected via PE
    static analysis. Sysmon registry_event (EID 13) and WMI activity logs
    are the primary compensating controls.

    EDR gap context (api_call_matrix.ps1):
      IWbemServices_ExecQuery (UH_WMI): COM vtable dispatch — all 7 EDRs
      rely on behavioral detection only, no direct API hook possible.
references:
    - https://attack.mitre.org/techniques/T1546/003/
    - https://attack.mitre.org/techniques/T1047/
    - https://www.ired.team/offensive-security/persistence/t1084-windows-management-instrumentation-event-subscription
author: Loaded Potato / lasersharkkiller
date: 2026-04-07
tags:
    - attack.persistence
    - attack.t1546.003
    - attack.defense_evasion
    - attack.t1047
logsource:
    category: registry_event
    product: windows
detection:
    selection_wmi_consumer:
        TargetObject|contains:
            - '\__EventConsumer'
            - 'ActiveScriptEventConsumer'
            - 'CommandLineEventConsumer'
    selection_suspicious_command:
        Details|contains:
            - 'powershell'
            - 'cmd.exe'
            - 'wscript'
            - 'cscript'
            - 'rundll32'
            - 'certutil'
            - 'bitsadmin'
            - 'mshta'
            - 'regsvr32'
            - 'installutil'
    condition: selection_wmi_consumer and selection_suspicious_command
falsepositives:
    - SCCM/ConfigMgr and some AV/EDR products legitimately create WMI subscriptions
    - Investigate the CommandLineTemplate or Script value to distinguish legitimate from malicious
level: high
