# detections/sigma/custom/PicusSim2026_T1053_MaliciousScheduledTask_StagingPath.yml
title: Scheduled Task Created Referencing Suspicious Path or Masqueraded Name (Picus Sim 2026-04-12)
id: d4e5f6a7-2010-4bbb-c010-picus20260412
status: experimental
description: >
    Detects schtasks.exe creating tasks that reference executables in user-writable
    staging paths (AppData\Roaming, Public) with names masquerading as legitimate
    software (Chrome.exe, GoogleUpdate.exe, GoogleCrash.exe). Also detects tasks
    using PowerShell bypass execution or tasks scheduled at very frequent intervals
    (every 3-10 minutes) typical of C2 beaconing persistence. Observed in Picus
    simulation 2026-04-12 with three separate tasks using masqueraded Chrome/Google
    names in AppData\Roaming staging paths.
    T1053.005 - Scheduled Task/Job: Scheduled Task.
references:
    - https://attack.mitre.org/techniques/T1053/005/
author: Loaded Potato - Picus Simulation 2026-04-12
date: 2026-04-13
tags:
    - attack.persistence
    - attack.execution
    - attack.t1053.005
    - picus.simulation.2026-04-12
logsource:
    category: process_creation
    product: windows
detection:
    selection_schtasks_create:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains: '/create'
    filter_masqueraded_staging:
        CommandLine|contains:
            - '\AppData\Roaming\Chrome\'
            - '\AppData\Roaming\Gchrome\'
            - '\AppData\Roaming\Google\'
            - '\Users\Public\Publ'
    filter_bypass_ps:
        CommandLine|contains|all:
            - 'powershell'
            - '-ep bypass'
            - '\Temp\'
    filter_high_frequency:
        CommandLine|re: '/MO\s+(0?[1-9]|10)\s'
        CommandLine|contains:
            - '/SC MINUTE'
    condition: selection_schtasks_create and (filter_masqueraded_staging or filter_bypass_ps or filter_high_frequency)
falsepositives:
    - Legitimate software installers creating frequent scheduled tasks (should be investigated)
level: high
