# Hardening/LinuxHardening/Configs/cis-auditd.rules
# ============================================================
# CIS Distribution Independent Linux Benchmark - auditd Rules
# Reference: CIS DIL Benchmark v2.0 / NIST 800-53 Rev 5 AU-12
# Generated by: LinuxHardening module (Set-AuditdPolicy)
#
# IMPORTANT: This file is deployed to:
#   /etc/audit/rules.d/99-hardened.rules
#
# Set-AuditdPolicy also appends rules for SUID/SGID binaries
# discovered dynamically on the target system.
# ============================================================

# Delete all existing rules first
-D

# Set the audit backlog buffer size (increase for high-traffic systems)
-b 8192

# Failure mode: 1 = printk (log and continue), 2 = panic (halt on failure)
# WARNING: Use -f 2 only if audit continuity is mandatory; it will halt the system
-f 1

# ============================================================
# CIS 4.1.3 - Date/time modification events (NIST AU-12)
# ============================================================
-a always,exit -F arch=b64 -S adjtimex,settimeofday -k time-change
-a always,exit -F arch=b32 -S adjtimex,settimeofday,stime -k time-change
-a always,exit -F arch=b64 -S clock_settime -k time-change
-a always,exit -F arch=b32 -S clock_settime -k time-change
-w /etc/localtime -p wa -k time-change

# ============================================================
# CIS 4.1.4 - User/group information modification (NIST AU-12)
# ============================================================
-w /etc/group    -p wa -k identity
-w /etc/passwd   -p wa -k identity
-w /etc/gshadow  -p wa -k identity
-w /etc/shadow   -p wa -k identity
-w /etc/security/opasswd -p wa -k identity

# ============================================================
# CIS 4.1.5 - Network environment changes (NIST AU-12)
# ============================================================
-a always,exit -F arch=b64 -S sethostname,setdomainname -k system-locale
-a always,exit -F arch=b32 -S sethostname,setdomainname -k system-locale
-w /etc/issue          -p wa -k system-locale
-w /etc/issue.net      -p wa -k system-locale
-w /etc/hosts          -p wa -k system-locale
-w /etc/hostname       -p wa -k system-locale
# RHEL/CentOS/Fedora
-w /etc/sysconfig/network          -p wa -k system-locale
-w /etc/sysconfig/network-scripts/ -p wa -k system-locale
# Debian/Ubuntu
-w /etc/network/interfaces -p wa -k system-locale
# Netplan (Ubuntu 18+)
-w /etc/netplan/ -p wa -k system-locale

# ============================================================
# CIS 4.1.6 - MAC policy changes (SELinux / AppArmor) (NIST AU-12)
# ============================================================
-w /etc/selinux/         -p wa -k MAC-policy
-w /usr/share/selinux/   -p wa -k MAC-policy
-w /etc/apparmor/        -p wa -k MAC-policy
-w /etc/apparmor.d/      -p wa -k MAC-policy

# ============================================================
# CIS 4.1.7 - Login and logout events (NIST AU-12)
# ============================================================
-w /var/log/lastlog       -p wa -k logins
-w /var/run/faillock/     -p wa -k logins
-w /var/log/tallylog      -p wa -k logins

# ============================================================
# CIS 4.1.8 - Session initiation (NIST AU-12)
# ============================================================
-w /var/run/utmp   -p wa -k session
-w /var/log/wtmp   -p wa -k logins
-w /var/log/btmp   -p wa -k logins

# ============================================================
# CIS 4.1.9 - DAC permission modification events (NIST AU-12)
# ============================================================
-a always,exit -F arch=b64 -S chmod,fchmod,fchmodat -F auid>=1000 -F auid!=4294967295 -k perm_mod
-a always,exit -F arch=b32 -S chmod,fchmod,fchmodat -F auid>=1000 -F auid!=4294967295 -k perm_mod
-a always,exit -F arch=b64 -S chown,fchown,fchownat,lchown -F auid>=1000 -F auid!=4294967295 -k perm_mod
-a always,exit -F arch=b32 -S chown,fchown,fchownat,lchown -F auid>=1000 -F auid!=4294967295 -k perm_mod
-a always,exit -F arch=b64 -S setxattr,lsetxattr,fsetxattr,removexattr,lremovexattr,fremovexattr -F auid>=1000 -F auid!=4294967295 -k perm_mod
-a always,exit -F arch=b32 -S setxattr,lsetxattr,fsetxattr,removexattr,lremovexattr,fremovexattr -F auid>=1000 -F auid!=4294967295 -k perm_mod

# ============================================================
# CIS 4.1.10 - Unsuccessful file access attempts (NIST AU-12)
# ============================================================
-a always,exit -F arch=b64 -S creat,open,openat,truncate,ftruncate -F exit=-EACCES -F auid>=1000 -F auid!=4294967295 -k access
-a always,exit -F arch=b32 -S creat,open,openat,truncate,ftruncate -F exit=-EACCES -F auid>=1000 -F auid!=4294967295 -k access
-a always,exit -F arch=b64 -S creat,open,openat,truncate,ftruncate -F exit=-EPERM  -F auid>=1000 -F auid!=4294967295 -k access
-a always,exit -F arch=b32 -S creat,open,openat,truncate,ftruncate -F exit=-EPERM  -F auid>=1000 -F auid!=4294967295 -k access

# ============================================================
# CIS 4.1.11 - Privileged command execution (NIST AU-12)
# Static entries; Set-AuditdPolicy appends SUID/SGID binaries
# discovered dynamically on the target system.
# ============================================================
-a always,exit -F path=/usr/bin/sudo       -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/su         -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/newgrp     -F perm=x -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/chsh       -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/chfn       -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/passwd     -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/gpasswd    -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/bin/crontab    -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/sbin/usermod   -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/sbin/useradd   -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/sbin/userdel   -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/usr/sbin/groupmod  -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged
-a always,exit -F path=/sbin/unix_chkpwd   -F perm=x -F auid>=1000 -F auid!=4294967295 -k privileged

# ============================================================
# CIS 4.1.12 - File system mounts (NIST AU-12)
# ============================================================
-a always,exit -F arch=b64 -S mount -F auid>=1000 -F auid!=4294967295 -k mounts
-a always,exit -F arch=b32 -S mount -F auid>=1000 -F auid!=4294967295 -k mounts

# ============================================================
# CIS 4.1.13 - File deletion events (NIST AU-12)
# ============================================================
-a always,exit -F arch=b64 -S unlink,unlinkat,rename,renameat -F auid>=1000 -F auid!=4294967295 -k delete
-a always,exit -F arch=b32 -S unlink,unlinkat,rename,renameat -F auid>=1000 -F auid!=4294967295 -k delete

# ============================================================
# CIS 4.1.14 - Sudoers changes (NIST AU-12, CM-3)
# ============================================================
-w /etc/sudoers    -p wa -k scope
-w /etc/sudoers.d/ -p wa -k scope

# ============================================================
# CIS 4.1.15 - sudo log (NIST AU-12)
# ============================================================
-w /var/log/sudo.log -p wa -k actions

# ============================================================
# CIS 4.1.16 - Kernel module load/unload (NIST AU-12, CM-6)
# ============================================================
-w /sbin/insmod      -p x -k modules
-w /sbin/rmmod       -p x -k modules
-w /sbin/modprobe    -p x -k modules
-w /usr/sbin/insmod  -p x -k modules
-w /usr/sbin/rmmod   -p x -k modules
-w /usr/sbin/modprobe -p x -k modules
-a always,exit -F arch=b64 -S init_module,delete_module -k modules
-a always,exit -F arch=b32 -S init_module,delete_module -k modules

# ============================================================
# Additional: Audit config and tool changes (NIST AU-9)
# ============================================================
-w /etc/audit/        -p wa -k auditconfig
-w /etc/libaudit.conf -p wa -k auditconfig
-w /etc/audisp/       -p wa -k audispconfig
-w /sbin/auditctl     -p x  -k audittools
-w /sbin/auditd       -p x  -k audittools

# ============================================================
# Additional: Cron job modifications (NIST CM-3)
# ============================================================
-w /etc/cron.allow     -p wa -k cron
-w /etc/cron.deny      -p wa -k cron
-w /etc/cron.d/        -p wa -k cron
-w /etc/cron.daily/    -p wa -k cron
-w /etc/cron.hourly/   -p wa -k cron
-w /etc/cron.monthly/  -p wa -k cron
-w /etc/cron.weekly/   -p wa -k cron
-w /etc/crontab        -p wa -k cron
-w /var/spool/cron/crontabs/ -p wa -k cron

# ============================================================
# Additional: SSH configuration changes (NIST CM-6, AC-17)
# ============================================================
-w /etc/ssh/sshd_config -p warx -k sshd
-w /etc/ssh/             -p wa   -k sshd

# ============================================================
# Additional: PAM configuration (NIST IA-3, CM-6)
# ============================================================
-w /etc/pam.d/                    -p wa -k pam
-w /etc/security/limits.conf      -p wa -k pam
-w /etc/security/pam_env.conf     -p wa -k pam
-w /etc/security/namespace.conf   -p wa -k pam
-w /etc/security/namespace.init   -p wa -k pam

# ============================================================
# Additional: Sysctl configuration (NIST CM-6)
# ============================================================
-w /etc/sysctl.conf      -p wa -k sysctl
-w /etc/sysctl.d/        -p wa -k sysctl

# ============================================================
# CIS 4.1.17 - Make audit configuration immutable (NIST AU-9)
#
# WARNING: Once set, audit rules cannot be changed without rebooting.
# Comment this line out or use -NoImmutable if you need to modify
# rules without a reboot during testing.
# ============================================================
-e 2
