# detections/tanium-picus/TR-PICUS-PERSIST-001_Scheduled_Task_Creation_PowerShell.json
{
    "type": "tanium-signal",
    "typeVersion": "1.0",
    "isSchemaValid": true,
    "createdAt": "2026-04-14T00:00:00.000Z",
    "updatedAt": "2026-04-14T00:00:00.000Z",
    "data": {
        "id": "TR-PICUS-PERSIST-001: Scheduled Task Creation via PowerShell Register-ScheduledTask",
        "name": "TR-PICUS-PERSIST-001: Scheduled Task Creation via PowerShell Register-ScheduledTask",
        "description": "Detects scheduled task creation via PowerShell's Register-ScheduledTask cmdlet or New-ScheduledTask. This is a common persistence mechanism used by attackers. Monitoring for this behavior helps detect post-exploitation persistence techniques. MITRE T1053.005.",
        "contents": "(process.name contains 'powershell.exe' OR process.name contains 'pwsh.exe') AND (process.command_line contains 'Register-ScheduledTask' OR process.command_line contains 'New-ScheduledTask' OR process.command_line contains 'New-ScheduledTaskAction' OR process.command_line contains 'New-ScheduledTaskTrigger' OR process.command_line contains 'schtasks.exe /create')",
        "syntax_version": 1,
        "mitreAttack": {
            "techniques": [
                {
                    "id": "T1053.005",
                    "name": "Scheduled Task/Job: Scheduled Task"
                },
                {
                    "id": "T1059.001",
                    "name": "Command and Scripting Interpreter: PowerShell"
                }
            ]
        },
        "platforms": ["windows"],
        "labels": ["Windows", "Persistence", "Picus", "PowerShell"]
    }
}
