# detections/tanium-picus/TR-PICUS-TIER1-SCHEDULED_TASK_POWERSHELL.json
{
    "type": "tanium-signal",
    "typeVersion": "1.0",
    "isSchemaValid": true,
    "createdAt": "2026-04-14T00:00:00.000Z",
    "updatedAt": "2026-04-14T00:00:00.000Z",
    "data": {
        "id": "TR-PICUS-TIER1-SCHEDULED_TASK_POWERSHELL: PowerShell Register-ScheduledTask Persistence",
        "name": "TR-PICUS-TIER1-SCHEDULED_TASK_POWERSHELL: PowerShell Register-ScheduledTask Persistence Creation",
        "description": "CRITICAL PERSISTENCE DETECTION: Detects scheduled task creation via PowerShell cmdlets, a primary persistence mechanism observed in 18+ undetected Picus instances. Monitors for PowerShell Register-ScheduledTask, New-ScheduledTask, New-ScheduledTaskAction, and New-ScheduledTaskTrigger cmdlet execution. Unlike schtasks.exe (which is already detected), PowerShell-based task creation often evades detection because alerts focus on the binary, not the cmdlet. Observable from detonation: powershell.exe with Register-ScheduledTask parameters, often from non-administrator context. Critical for post-exploitation persistence.",
        "contents": "(process.name contains 'powershell.exe' OR process.name contains 'pwsh.exe') AND (process.command_line contains 'Register-ScheduledTask' OR process.command_line contains 'New-ScheduledTask' OR process.command_line contains 'New-ScheduledTaskAction' OR process.command_line contains 'New-ScheduledTaskTrigger' OR process.command_line contains 'New-ScheduledTaskPrincipal' OR process.command_line contains 'Set-ScheduledTask') AND (process.parent.name != 'explorer.exe' AND process.parent.name != 'svchost.exe')",
        "syntax_version": 1,
        "mitreAttack": {
            "techniques": [
                {
                    "id": "T1053.005",
                    "name": "Scheduled Task/Job: Scheduled Task"
                },
                {
                    "id": "T1059.001",
                    "name": "Command and Scripting Interpreter: PowerShell"
                }
            ]
        },
        "platforms": ["windows"],
        "labels": ["Windows", "Persistence", "Picus", "PowerShell", "Scheduled Tasks", "18+ Instances"],
        "expectedCoverage": "18+ undetected scheduled task persistence instances"
    }
}
